# sa/⭕ 当一条合格的音游狗！/daemon.conf
[program:watchman]
directory = "SOMEWHERE"
command = "/bin/env python3 -m SOMEWHERE/watchman"
